# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=service nginx start
autostart=true
autorestart=true

[program:php7.1-fpm]
command=service php7.1-fpm start
autostart=true
autorestart=true